# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:api]
directory=/app/apps/api
command=/bin/bash -lc 'uv run python -m api.main'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=APP_ENV="prod",APP_HOST="127.0.0.1",APP_PORT="8001"

[program:web]
directory=/app/apps/web
command=/bin/bash -lc 'corepack enable >/dev/null 2>&1 || true && yarn start'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
